# systemd/sysreport.service
[Unit]
Description=System Health Monitor Service

[Service]
Type=oneshot
ExecStart=__REPO_DIR__/scripts/sysreport.sh

[Install]
WantedBy=multi-user.target
